# supervisord.conf
[supervisord]
nodaemon=true

[program:WowzaStreamingEngine]
command=/usr/local/bin/wowza-engine-start

[program:WowzaStreamingMananger]
autostart=false
command=/usr/local/bin/wowza-manager-start
